# etc/back-up-postgresql.timer
[Unit]
Description=Back up PostgreSQL

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
